# systemd/network-manager.service
[Unit]
Description=Network Manager: Detects connectivity and manages WiFi hotspot
Wants=network-online.target
After=network-online.target
Before=stream-player.service
# Ensure WiFi interface is available
After=sys-subsystem-net-devices-wlan0.device

[Service]
Type=simple
# Run as root to manage network interfaces
User=root
Environment=NETWORK_WAIT_TIMEOUT=30
Environment=HOTSPOT_SSID=bartix-config
Environment=HOTSPOT_PASSWORD=bartix-config
Environment=HOTSPOT_INTERFACE=wlan0
Environment=HOTSPOT_IP=192.168.4.1
# Don't wait too long for network-online (it may never come)
TimeoutStartSec=60
# Pre-start: Ensure WiFi is unblocked and ready (if interface exists)
ExecStartPre=/bin/sh -c 'if [ -e /sys/class/net/wlan0 ]; then rfkill unblock wifi || true; sleep 1; iw dev wlan0 set txpower fixed 2000 || true; fi'
ExecStart=/usr/bin/env python3 /usr/local/bin/network-manager.py
Restart=always
RestartSec=5

# Resource limits
LimitNOFILE=1024

# Ensure we can manage network interfaces
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target

